# useful/systemdServices/example.service
[Unit]
Description=example service
[Service]
WorkingDirectory=/vagrant/
User=root
ExecStart= /usr/bin/node index.js
Restart=always

